# archive/LaunchDaemons-System/us.t0m.monthlyZFSscrubOfTank.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>us.t0m.monthlyZFSscrubOfTank</string>

    <key>ProgramArguments</key>
    <array>
      <string>/usr/local/bin/zpool</string>
      <string>scrub</string>
      <string>tank</string>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <key>StartCalendarInterval</key>
    <array>
      <dict>
        <key>Day</key>
        <integer>1</integer>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>1</integer>
        <key>Minute</key>
        <integer>11</integer>
      </dict>
      <dict>
        <key>Day</key>
        <integer>2</integer>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>1</integer>
        <key>Minute</key>
        <integer>11</integer>
      </dict>
      <dict>
        <key>Day</key>
        <integer>3</integer>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>1</integer>
        <key>Minute</key>
        <integer>11</integer>
      </dict>
      <dict>
        <key>Day</key>
        <integer>4</integer>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>1</integer>
        <key>Minute</key>
        <integer>11</integer>
      </dict>
      <dict>
        <key>Day</key>
        <integer>5</integer>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>1</integer>
        <key>Minute</key>
        <integer>11</integer>
      </dict>
      <dict>
        <key>Day</key>
        <integer>6</integer>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>1</integer>
        <key>Minute</key>
        <integer>11</integer>
      </dict>
      <dict>
        <key>Day</key>
        <integer>7</integer>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>1</integer>
        <key>Minute</key>
        <integer>11</integer>
      </dict>
    </array>
  </dict>
</plist>
