# archive/LaunchDaemons-System/us.t0m.syncZFSsnapshots.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>us.t0m.syncZFSsnapshots</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/tom/bin/syncZFSsnapshots.sh</string>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <key>StartInterval</key>
    <integer>3600</integer>

    <key>StandardErrorPath</key>
    <string>/tmp/syncZFSsnapshots.err</string>

    <key>StandardOutPath</key>
    <string>/tmp/syncZFSsnapshots.out</string>
</dict>
</plist>
